# dist/profile.d/chromatindb.sh
# Canonical paths for chromatindb CLI tools.
# Sourced automatically by login shells via /etc/profile.d/.
# Peer-management subcommands (add-peer, remove-peer, list-peers)
# consult this variable first when --config is not provided.
export CHROMATINDB_CONFIG=/etc/chromatindb/node.json
